# etc/yai/profile.d/yai.sh
# shellcheck shell=sh
export YAI_HOME="${YAI_HOME:-$HOME/.yai}"
export YAI_RUNTIME_DIR="${YAI_RUNTIME_DIR:-/run/yai}"
export YAI_STATE_DIR="${YAI_STATE_DIR:-/var/lib/yai}"
export YAI_CACHE_DIR="${YAI_CACHE_DIR:-/var/cache/yai}"
export YAI_LOG_DIR="${YAI_LOG_DIR:-/var/log/yai}"
export YAI_TMP_DIR="${YAI_TMP_DIR:-/tmp/yai}"
