# usr/lib/systemd/user/yai-sh.service
[Unit]
Description=YAI Interactive Shell Session

[Service]
Type=simple
ExecStart=/usr/bin/yai-sh
Restart=no

[Install]
WantedBy=default.target
